# .config/systemd/user/change-wallpaper.timer
[Unit]
Description=Run timer to change wallpaper

[Timer]
OnActiveSec=0
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
