# .config/systemd/user/change-alacritty-theme.service
[Unit]
Description=Change Alacritty Theme

[Service]
ExecStart=%h/.local/bin/change-alacritty-theme -t d

[Install]
WantedBy=default.target
