# .config/systemd/user/network-check.service
[Unit]
Description=Monitor Network

[Service]
Type=oneshot
ExecStart=nc -z www.google.com 80
