# .config/systemd/user/network-check.timer
[Unit]
Description=Run timer to check network connection

[Timer]
OnActiveSec=0
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
